# dehydrated-alpine/crontab
41 6,15 * * * /usr/local/bin/dehydrated -c